# rootfs/etc/local.d/expand-rootfs.start
#!/bin/sh
/usr/local/bin/expand-rootfs.sh >> /var/log/expand-rootfs.log 2>&1
